# sensu_configs/init/sensu-dashboard.conf
# Upstart: /etc/init/sensu-dashboard.conf

description "sensu dashboard"

start on (local-filesystems and net-device-up)
stop on [!12345]

respawn

chdir /opt/sensu/bin
env PATH=/etc/sensu/plugins:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

pre-start script
  test -x ./sensu-dashboard || { stop; exit 0; }
end script

exec setuidgid sensu ./sensu-dashboard -l /var/log/sensu/sensu-dashboard.log
